# systemd/polymarket-bot.service
[Unit]
Description=Polymarket YES/NO Arbitrage Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polymarket
WorkingDirectory=/opt/polymarket-bot
EnvironmentFile=/opt/polymarket-bot/.env
ExecStart=/usr/bin/python3 run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polymarket-bot

[Install]
WantedBy=multi-user.target
